# dist/gadael.service
[Unit]
Description=Gadael HTTP server
After=network.target

[Service]
User=root
WorkingDirectory=/var/lib/gadael
ExecStart=/usr/bin/node app.js
StandardOutput=syslog

[Install]
WantedBy=multi-user.target
